# deploy/manufacturing-audit.timer
[Unit]
Description=Daily internal audit certification for the manufacturing plant

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
